# dot_config/systemd/user/gitea_copy_backup.timer
[Unit]
Description=Daily backup of Gitea on DigitalOcean droplet

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
