# dot_config/systemd/user/grocy_copy_backup.timer
[Unit]
Description=Daily backup of Grocy on DigitalOcean droplet

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
